# systemd/minecraft@.socket
[Unit]
PartOf=minecraft@%i.service

[Socket]
ListenFIFO=%t/minecraft@%i.stdin
